# usr/lib/opkg/info/elink.list
/usr/bin/elink
/etc/init.d/elink
/usr/sbin/elink_ap_realdev.lua
