# usr/lib/opkg/info/mi_iptv.list
/usr/sbin/mi_iptv.sh
/etc/config/mi_iptv
/etc/init.d/mi_iptv
